# modules/pz_base/files/connectivity-check.cron
# m h dom mon dow  user       command
*/15   *   *   *   root       /usr/local/bin/connectivity-check

