# modules/pz_sentinel/files/sentinel.cron
# m h   dom mon dow user     command
  0 1    *   *   *  sentinel sh -c 'cd /home/sentinel/sentinel && python3 /home/sentinel/sentinel/scheduled_tasks.py'

